# scripts/systemd/user/edgerun-cloudflared-term.service
[Unit]
Description=Edgerun Cloudflared Tunnel (term.edgerun.tech)
After=network-online.target edgerun-term-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
EnvironmentFile=-__HOME__/.config/edgerun/cloudflared-term.env
ExecStart=__ROOT_DIR__/scripts/cloudflared/run-term-tunnel.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__ROOT_DIR__ __HOME__/.cloudflared __HOME__/.config/edgerun

[Install]
WantedBy=default.target
