# scripts/systemd/user/edgerun-term-server.service
[Unit]
Description=Edgerun Terminal Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
EnvironmentFile=-__HOME__/.config/edgerun/term-server.env
ExecStart=__ROOT_DIR__/out/target/release/edgerun-term-server
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__ROOT_DIR__/out __HOME__/.config/edgerun

[Install]
WantedBy=default.target
