# scripts/systemd/user/solana-test-validator.service
[Unit]
Description=Local Solana Test Validator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ROOT_DIR__
ExecStart=/usr/bin/solana-test-validator \
  --ledger __HOME__/.config/edgerun/solana/ledger \
  --bind-address 127.0.0.1 \
  --rpc-port 8899 \
  --dynamic-port-range 9000-10000 \
  --limit-ledger-size 50000000
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__HOME__/.config/edgerun/solana

[Install]
WantedBy=default.target
